# crontab/backup-tularosa
0   5   *   *   sun  ~/bin/backups/to_tularosa >/dev/null
