# crontab/backup-tularosa
0   *   *   *   *  ~/bin/backups/to_tularosa >/dev/null
