# bin/crontab/autocommit
0   *   *   *   *  ~/etc/bin/autocommit.git ~/src/org-files >/dev/null
