# bin/crontab/backup-tularosa
0   5   *   *   sun  ~/etc/bin/backups/to_tularosa >/dev/null
